# package/DEBIAN/postinst
#!/bin/sh
systemctl enable device-init.service
